# .docker/supervisor/supervisord-batch-stg.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
user=www-data
directory=/var/www/html
command=
autostart=true
autorestart=true
numprocs=5
startretries=10
stdout_events_enabled=1
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
